# launchd/com.toss-trading.discord-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.toss-trading.discord-bot</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PATH</key>
        <string>/Users/aiden/Desktop/Auto-trader/tossinvest-cli/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>TOSSCTL_AUTH_HELPER_DIR</key>
        <string>/Users/aiden/Desktop/Auto-trader/tossinvest-cli/auth-helper</string>
        <key>TOSSCTL_AUTH_HELPER_PYTHON</key>
        <string>/Users/aiden/Desktop/Auto-trader/tossinvest-cli/auth-helper/.venv/bin/python3</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/aiden/Desktop/Auto-trader/toss-trading-system/.venv/bin/python3</string>
        <string>/Users/aiden/Desktop/Auto-trader/toss-trading-system/scripts/discord-bot.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/aiden/Desktop/Auto-trader/toss-trading-system</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/aiden/Library/Logs/toss-trading/discord-bot.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/aiden/Library/Logs/toss-trading/discord-bot.error.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
